# claude-irc-mcp.service
[Unit]
Description=2600net IRC MCP Server
Documentation=https://github.com/astrutt/claude-irc-mcp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claudemcp
Group=claudemcp
ExecStart=/opt/claude-irc-mcp/venv/bin/python3 /opt/claude-irc-mcp/irc_mcp_server.py
WorkingDirectory=/opt/claude-irc-mcp
Restart=on-failure
RestartSec=15s
StartLimitIntervalSec=300
StartLimitBurst=5
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-irc-mcp
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
